# tmpl/install/centreon.cron
#!/bin/sh
#
# Cron Configuration for Centreon
#
# rewrite file with new cron line
CRONTAB_EXEC_USER=""

# Cron for log analyser - for reporting
0 1 1-31 * * @CRONUSER@ @PHP_BIN@ -q @INSTALL_DIR_CENTREON@/cron/reporting/ArchiveLogInDB.php >> @CENTREON_LOG@/archiveLogInDB.log 2>&1

# Cron for CentACL
*/2 * * * * @CRONUSER@ @PHP_BIN@ -q @INSTALL_DIR_CENTREON@/cron/centAcl.php >> @CENTREON_LOG@/centAcl.log 2>&1
