# ci/debian/centreon-common.postrm
#!/bin/sh

set -e

case "$1" in
    remove)
        deluser centreon || true
        delgroup centreon || true
        systemctl --no-reload disable --now centreon.service || : &>/dev/null || :
    ;;
esac

exit 0
